# koji-offsite-backup/osg-koji-offsite-backup.cron
# place in /etc/cron.d:
0 6 * * sat root /root/osg-koji-offsite-backup
